# deployment/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=hypercorn backend.app:app --bind 127.0.0.1:5678
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
environment=PYTHONPATH="/app:/app/backend",PYTHONUNBUFFERED="1",TASKLINE_DEBUG="0"

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=10

[program:db-migration]
command=sh -c "cd /app/backend && alembic upgrade head"
directory=/app
autostart=true
autorestart=false
startsecs=0
stderr_logfile=/var/log/supervisor/db-migration.err.log
stdout_logfile=/var/log/supervisor/db-migration.out.log
priority=1
